# at-3.1.16/atd.service.in
[Unit]
Description=Deferred execution scheduler
Documentation=man:atd(8)

[Service]
ExecStart=@sbindir@/atd -f
IgnoreSIGPIPE=false

[Install]
WantedBy=multi-user.target
